# crontab
# min	hour	day	    month	weekday	command
*/30    *       *       *       *       /usr/local/bin/python /usr/src/app/antibox.py
