# crontab
*/1 * * * * root /usr/local/bin/python /usr/src/app/antibox.py
